# rtl4332mqtt/rootfs/etc/services.d/rtl4332mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# SmartHouse Hass.io Add-ons: RTL4332MQTT
# Runs the RTL433 TO MQTT
# ==============================================================================

# Wait for RTL-SDR to become available
bashio::net.wait_for 1234 localhost 300

# Run the RTL4332MQTT
bashio::log.info "Starting RTL433 2 MQTT..."
exec sh /app/rtl2mqtt.sh
